# src/garage-controller.service
[Unit]
Description=Garage Door Relay Controller and DHT22 Publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=caseman
EnvironmentFile=/opt/garage-controller/config.env
EnvironmentFile=/opt/garage-controller/secrets.env
ExecStart=/opt/garage-controller/venv/bin/python /opt/garage-controller/relay_controller.py
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
